# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:goapp]
command=/go/bin/main
autostart=false
autorestart=false
